# pi/midi-bridge.service
[Unit]
Description=Piano MIDI to WebSocket Bridge
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=finlay
ExecStartPre=/home/finlay/midi-bridge/deploy.sh
WorkingDirectory=/home/finlay/midi-bridge
ExecStart=/home/finlay/midi-bridge/venv/bin/python3 /home/finlay/midi-bridge/midi-bridge.py --www /home/finlay/midi-bridge/www --data /home/finlay/midi-bridge/data
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/finlay
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
